# etc/rc.d/openhapd
#!/bin/ksh
# /etc/rc.d/openhapd
#
# OpenHAP daemon rc.d script for OpenBSD

daemon="/usr/local/bin/openhapd"
daemon_flags="-c /etc/openhapd.conf"
# Note: daemon starts as root and drops to _openhap after privileged operations

. /etc/rc.d/rc.subr

rc_configtest() {
    ${daemon} -n ${daemon_flags}
}

rc_reload() {
    pkill -HUP -xf "${pexp}"
}

pexp="perl ${daemon}.*"

rc_cmd $1
